# docker/php/crontab/root
* * * * * php /var/www/laravel.su/artisan schedule:run > ./storage/logs/schedule.log 2>&1 &
